# overlay/etc/supervisor.d/dbus.ini
[program:dbus]
priority=10
command=/usr/bin/start-dbus
