# overlay/etc/supervisor.d/desktop.ini
[program:desktop]
priority=50
user=player
directory=/home/player
command=/usr/bin/start-desktop
stopasgroup=true
killasgroup=true
