# overlay/etc/supervisor.d/pipewire-pulse.ini
[program:pipewire-pulse]
priority=42
user=player
directory=/home/player
command=/usr/bin/start-pipewire-pulse
