# overlay/etc/supervisor.d/pipewire.ini
[program:pipewire]
priority=40
user=player
directory=/home/player
command=/usr/bin/start-pipewire
