# overlay/etc/supervisor.d/polkit.ini
[program:polkit]
priority=11
command=/usr/bin/start-polkit
