# overlay/etc/supervisor.d/seatd.ini
[program:seatd]
priority=9
command=/usr/bin/seatd -u player
environment=SEATD_VTBOUND="0"
