# overlay/etc/supervisor.d/session-dbus.ini
[program:session-dbus]
priority=20
user=player
directory=/home/player
command=/usr/bin/start-session-dbus
