# overlay/etc/supervisor.d/sunshine.ini
[program:sunshine]
priority=80
user=player
directory=/home/player
command=/usr/bin/start-sunshine
stopasgroup=true
killasgroup=true
