# overlay/etc/supervisor.d/udev.ini
[program:udev]
priority=10
command=/usr/bin/dumb-udev
