# overlay/etc/supervisor.d/wireplumber.ini
[program:wireplumber]
priority=41
user=player
directory=/home/player
command=/usr/bin/start-wireplumber
